# deployment/systemd/can-gateway.service
[Unit]
Description=Greenhouse CAN Gateway
After=network-online.target can0-up.service
Wants=network-online.target
Requires=can0-up.service
RequiresMountsFor=/data

[Service]
ExecStart=/usr/bin/can-gateway --config /etc/greenhouse/gateway.conf --interface can0
Restart=always
RestartSec=5
User=greenhouse
Group=greenhouse

[Install]
WantedBy=multi-user.target
