# deployment/systemd/cloud-sync.service
[Unit]
Description=Greenhouse Cloud Sync (gRPC client)
After=network-online.target
Wants=network-online.target
RequiresMountsFor=/data

[Service]
ExecStart=/usr/bin/cloud-sync --config /etc/greenhouse/gateway.conf
Restart=always
RestartSec=5
User=greenhouse
Group=greenhouse

[Install]
WantedBy=multi-user.target
